# RaspberryPiFiles/airodump.service
[Unit]
Description=airodump

[Service]
Type=idle
ExecStart=/home/pi/airodump.sh
